# debian/config-testbox.prerm
#!/bin/sh
set -e

PKG=config-testbox

if [ "$1" = remove ] ; then
    for f in hostname ser2net.conf
    do
	[ -L /etc/$f ] && rm /etc/$f
	dpkg-divert --remove --package ${PKG} --rename \
		    --divert /etc/$f.distrib /etc/$f
    done
fi

invoke-rc.d ser2net stop

#DEBHELPER#

exit 0
