# debian/config-testbox.postrm
#!/bin/sh
set -e

PKG=config-testbox

if [ "$1" = remove ] ; then
    rm -f /etc/testbox/hostname
    rm -f /etc/testbox/ser2net.conf
fi

invoke-rc.d ser2net start || true

#DEBHELPER#

exit 0
